# systemd/system/reportd.service.in
[Unit]
Description=ABRT Report Daemon

[Service]
Type=dbus
BusName=org.freedesktop.reportd
ExecStart=@libexecdir@/reportd --system
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictRealtime=true
