# systemd/user/reportd.service.in
[Unit]
Description=ABRT Report Daemon

[Service]
Type=dbus
BusName=org.freedesktop.reportd
ExecStart=@libexecdir@/reportd
PrivateTmp=true
